# legion-go-setup.service
[Unit]
Description=Legion Go Permission Setup
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/legion-go-setup.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
